# systemd/osbuild-composer-api.socket
[Unit]
Description=OSBuild Composer API socket

[Socket]
Service=osbuild-composer.service
ListenStream=10443

[Install]
WantedBy=sockets.target
